# assets/worker-phpfpm.conf
[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true
